# services/bestool-alertd.service
[Unit]
Description=Tamanu alert daemon
After=network.target

[Service]
# notify-reload: the daemon signals readiness (READY=1), status, and shutdown
# (STOPPING=1) over $NOTIFY_SOCKET, and `systemctl reload` sends SIGHUP (the
# default ReloadSignal) which the daemon answers with RELOADING=1 / READY=1 —
# re-registering backup capabilities (picking up /etc/bestool/backups changes)
# without a restart. Needs systemd >= 253; on older systemd use Type=notify plus
# `ExecReload=/bin/kill -HUP $MAINPID`.
Type=notify-reload
# Distinct journal identifier so rsyslog can filter alertd apart from other
# bestool invocations (which all otherwise log as "bestool").
SyslogIdentifier=bestool-alertd
ExecStart=/usr/bin/bestool alertd run
MemoryMax=500M
Restart=always
RestartSec=10

# Talk to the rootful podman API socket instead of forking a privileged podman.
# Every podman call becomes a thin API client — the privileged container work
# happens in the daemon — so this unit can drop caps, namespaces, devices, and
# @mount below. Requires `systemctl enable --now podman.socket`. The CLI reads
# CONTAINER_HOST from the env automatically.
Environment=CONTAINER_HOST=unix:///run/podman/podman.sock

# Doctor's DB checks connect as a dedicated read-only Postgres role over the
# socket, authenticated by peer (pg_ident maps root -> that role) so there's no
# password to keep in sync with the app. The connection string is host-specific
# (it varies per tamanu install), so it is not baked into this unit: the
# deployment drops a /etc/systemd/system/bestool-alertd.service.d/*.conf
# override with `Environment=TAMANU_DATABASE_URL=...` on the hosts that need it.

# Runs as root with a restricted capability bounding set rather than the full
# root set. As uid 0 the kernel grants the whole bounding set into the effective
# set on exec (and the btrfs child below re-derives it the same way), so listing
# the caps here is enough and AmbientCapabilities stays empty.
#   - CAP_DAC_READ_SEARCH: lets root read other-owned, mode-0600 files the
#     doctor checks inspect (e.g. a service's TLS key) without also gaining DAC
#     write-override. (Kopia no longer needs this — it runs as the kopia user
#     and reads its own config.)
#   - CAP_SYS_ADMIN: the btrfs disk-stats doctor check shells out to `btrfs
#     device stats`, whose device-info ioctls return EPERM without it.
# CAP_SYS_ADMIN is broad (near root-equivalent), so the Protect*/seccomp/
# namespace directives below carry the hardening here, not the capability set.
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_SYS_ADMIN
AmbientCapabilities=
NoNewPrivileges=yes

ProtectSystem=strict
# read-only, NOT `yes`: a tamanu install may live under /home.
ProtectHome=read-only
CacheDirectory=bestool
# Point the `dirs` crate's cache at /var/cache/bestool (the CacheDirectory).
Environment=XDG_CACHE_HOME=%C
# Root-owned scratch for backups: the daemon (root, without DAC write-override)
# creates the postgres snapshot mounts and the base-backup staging here, under
# /var/lib/bestool, then hands them to the kopia user — it can't write into the
# kopia user's own home below.
StateDirectory=bestool
# bestool writes to its config dir, which ProtectSystem=strict otherwise
# locks read-only along with the rest of /etc.
ReadWritePaths=/etc/bestool
# The kopia user's home: its repo config (.config/kopia) and its cache and logs
# (.cache/kopia, since we run kopia as the kopia user with HOME here). kopia
# writes the cache, so this must be writable.
ReadWritePaths=/var/lib/kopia
# podman cp's tempdir lands in a private /tmp.
PrivateTmp=yes

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
# reads time-sync status, never sets the clock.
ProtectClock=yes
ProtectHostname=yes
# but NOT ProcSubset=pid — sysinfo needs /proc/meminfo and /proc/loadavg.
ProtectProc=invisible

# postgres snapshot backups (btrfs/thin-LVM) mount the snapshot block device
# (/dev/disk/by-uuid/*, /dev/<vg>/<lv>) read-only for kopia; a private /dev
# would hide it. The container work still goes to the remote rootful podman, so
# nothing else here needs a real device node.
PrivateDevices=no
# The read-only snapshot is mounted with X-mount.idmap so the kopia user can
# read postgres-owned files, and an idmapped mount creates a user namespace;
# this allows that one type and no other. The container work still goes to the
# remote rootful podman, so no other namespace is needed.
RestrictNamespaces=user
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
# @mount (mount/umount/mount_setattr/…): snapshot backups mount the snapshot
# read-only for kopia and unmount it afterwards; without it mount() is EPERM.
SystemCallFilter=@mount
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
# 0007, not 0077: /etc/bestool is a setgid tamanu-group dir shared with
# interactive (login user) and kopia/tamanu bestool invocations, so files this
# service writes there must stay group-readable and -writable. Other gets none.
UMask=0007

[Install]
WantedBy=multi-user.target
