# services/bestool-alertd.service
[Unit]
Description=Tamanu alert daemon
After=network.target

[Service]
# notify-reload: the daemon signals readiness (READY=1), status, and shutdown
# (STOPPING=1) over $NOTIFY_SOCKET, and `systemctl reload` sends SIGHUP (the
# default ReloadSignal) which the daemon answers with RELOADING=1 / READY=1 —
# re-registering backup capabilities (picking up /etc/bestool/backups changes)
# without a restart. Needs systemd >= 253; on older systemd use Type=notify plus
# `ExecReload=/bin/kill -HUP $MAINPID`.
Type=notify-reload
# Distinct journal identifier so rsyslog can filter alertd apart from other
# bestool invocations (which all otherwise log as "bestool").
SyslogIdentifier=bestool-alertd
ExecStart=/usr/bin/bestool alertd run
MemoryMax=500M
Restart=always
RestartSec=10

# Talk to the rootful podman API socket instead of forking a privileged podman.
# Every podman call becomes a thin API client — the privileged container work
# happens in the daemon — so this unit can drop caps, namespaces, devices, and
# @mount below. Requires `systemctl enable --now podman.socket`. The CLI reads
# CONTAINER_HOST from the env automatically.
Environment=CONTAINER_HOST=unix:///run/podman/podman.sock

# Doctor's DB checks connect as a dedicated read-only Postgres role over the
# socket, authenticated by peer (pg_ident maps root -> that role) so there's no
# password to keep in sync with the app. The connection string is host-specific
# (it varies per tamanu install), so it is not baked into this unit: the
# deployment drops a /etc/systemd/system/bestool-alertd.service.d/*.conf
# override with `Environment=TAMANU_DATABASE_URL=...` on the hosts that need it.

# Runs as root with a restricted capability bounding set rather than the full
# root set. As uid 0 the kernel grants the whole bounding set into the effective
# set on exec (and the btrfs child below re-derives it the same way), so listing
# the caps here is enough and AmbientCapabilities stays empty.
#   - CAP_DAC_READ_SEARCH: root reads any file (e.g. kopia's mode-0600 repo
#     config) without also gaining DAC write-override. Without it kopia flips to
#     wanting sudo, which NoNewPrivileges then blocks.
#   - CAP_SYS_ADMIN: the btrfs disk-stats doctor check shells out to `btrfs
#     device stats`, whose device-info ioctls return EPERM without it.
# CAP_SYS_ADMIN is broad (near root-equivalent), so the Protect*/seccomp/
# namespace directives below carry the hardening here, not the capability set.
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_SYS_ADMIN
AmbientCapabilities=
NoNewPrivileges=yes

ProtectSystem=strict
# read-only, NOT `yes`: a tamanu install may live under /home.
ProtectHome=read-only
# Two managed cache dirs under /var/cache (everything else there is read-only
# under ProtectSystem=strict): `bestool` for our own cache, and `kopia` because
# kopia derives its log/cache directory from XDG_CACHE_HOME (set to /var/cache
# below) as <XDG_CACHE_HOME>/kopia, i.e. /var/cache/kopia — without this it
# fails the backup with "mkdir /var/cache/kopia: read-only file system".
CacheDirectory=bestool kopia
# XDG_CACHE_HOME=/var/cache: the `dirs` crate appends "bestool" (→ the bestool
# CacheDirectory) and kopia appends "kopia" (→ the kopia CacheDirectory above).
Environment=XDG_CACHE_HOME=%C
# bestool writes to its config dir, which ProtectSystem=strict otherwise
# locks read-only along with the rest of /etc.
ReadWritePaths=/etc/bestool
# The canopy-managed kopia repo config lives at /var/lib/kopia/.config, and
# postgres backups stage their snapshot/base-backup under
# /var/lib/kopia/bestool-backup before kopia reads it — both need to be writable.
ReadWritePaths=/var/lib/kopia
# podman cp's tempdir lands in a private /tmp.
PrivateTmp=yes

ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
# reads time-sync status, never sets the clock.
ProtectClock=yes
ProtectHostname=yes
# but NOT ProcSubset=pid — sysinfo needs /proc/meminfo and /proc/loadavg.
ProtectProc=invisible

# no, not yes: postgres snapshot backups (btrfs/thin-LVM) mount the snapshot
# block device (/dev/disk/by-uuid/*, /dev/<vg>/<lv>) read-only for kopia, which
# a private /dev would hide. The container work still goes to the remote rootful
# podman, so no local device node is needed for that.
PrivateDevices=no
# user only (not yes = block all): the read-only snapshot is mounted with
# X-mount.idmap so the kopia user can read postgres-owned files, and an
# idmapped mount creates a user namespace. No other namespace type is allowed —
# the container work still goes to the remote rootful podman.
RestrictNamespaces=user
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
SystemCallFilter=@system-service
# @mount (mount/umount/mount_setattr/…): snapshot backups mount the snapshot
# read-only for kopia and unmount it afterwards; without it mount() is EPERM.
SystemCallFilter=@mount
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
# 0007, not 0077: /etc/bestool is a setgid tamanu-group dir shared with
# interactive (login user) and kopia/tamanu bestool invocations, so files this
# service writes there must stay group-readable and -writable. Other gets none.
UMask=0007

[Install]
WantedBy=multi-user.target
